# deploy/tunnel.service
[Unit]
Description=roboat tunnel service (connection layer)
After=network.target

[Service]
Type=simple
User=roboat
WorkingDirectory=/opt/roboat
ExecStart=/opt/roboat/bin/tunnel-svc
Restart=always
RestartSec=5
EnvironmentFile=/opt/roboat/config/.env
StandardOutput=append:/var/log/roboat.log
StandardError=append:/var/log/roboat.log

[Install]
WantedBy=multi-user.target
